# raw/etc/rc.local
#!/bin/sh -e
#
# rc.local
#

/usr/local/bin/minereu.sh
printf "Happy mining!\n"

exit 0

