# raw/etc/rc.local
#!/bin/sh -e
#
# rc.local
#

printf "Happy mining!\n"

exit 0

