# raw/etc/cron.d/minereu
10 *     * * *     root   /usr/sbin/ntpdate -u 0.uk.pool.ntp.org
*/2 *     * * *     root  /usr/local/bin/minereu.sh -o scan -p minereu_ -c uk -l bh -d p2 -s 0
* *     * * *     root  /usr/local/bin/worker_monitor.sh
